# ha_stats/rootfs/etc/services.d/ha-stats/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the HA Stats Lake service
# ==============================================================================
bashio::log.info "Starting HA Stats Lake..."
exec python3 /app/ha_stats.py
